# .config/systemd/user/dotfiles-backup.timer
[Unit]
Description=Daily dotfiles backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
