# .config/systemd/user/dotfiles-backup.service
[Unit]
Description=Regular dotfiles backup

[Service]
Type=oneshot
ExecStart=%h/bin/backup_dotfiles.sh
