# .config/systemd/user/libinput-gestures.service
[Unit]
Description=Libinput gestures daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# If your systemd --user already has DISPLAY/XAUTHORITY, you don't need these.
# If gestures work but key-sends fail, uncomment the next two lines:
#Environment=DISPLAY=:0
#Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/usr/bin/libinput-gestures
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target

