# bin/.config/systemd/user/copyq.timer
[Unit]
Description=CopyQ clipboard manager

[Service]
Type=simple
ExecStart=/usr/bin/copyq

[Install]
WantedBy=default.target
