# bin/.config/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon
After=default.target

[Service]
ExecStart=/usr/bin/dunst
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=xsession.target
