# bin/.config/systemd/user/picom.service
[Unit]
Description=Picom compositor
After=default.target

[Service]
ExecStart=%h/.local/bin/picom-autoconfig
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=xsession.target

