# bin/.config/systemd/user/picom.service
[Unit]
Description=Picom compositor
After=default.target

[Service]
ExecStart=/usr/bin/picom --config=%h/.config/picom/picom.conf
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=xsession.target

