# bin/.config/systemd/user/xsession.target
[Unit]
Description=User X session is ready
