# systemd-etc/etc/systemd/system/luks-backup.service
[Unit]
Description=LUKS Header Backup Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/luks-backup.sh
User=root
StandardOutput=journal
StandardError=journal
